# task-auto-update@.service
[Unit]
Description=Task Management System Auto-Update Monitor
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/University/7 season/PaperWork
ExecStart=/home/%i/University/7 season/PaperWork/auto-update.sh
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=task-auto-update

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
